# systemd/uu-remote-bridge.service
[Unit]
Description=UU Remote bridge for the live GNOME desktop
After=default.target

[Service]
Type=simple
EnvironmentFile=-%h/.config/uu-remote-bridge/environment
ExecStart=%h/.local/bin/uu-remote-bridge
Restart=on-failure
RestartSec=3
TimeoutStopSec=15
KillMode=control-group

[Install]
WantedBy=default.target
