# systemd/uu-remote-repair-monitor.timer
[Unit]
Description=Persistent UU repair and health monitor

[Timer]
OnBootSec=7min
OnUnitInactiveSec=15min
RandomizedDelaySec=2min
Unit=uu-remote-repair-monitor.service

[Install]
WantedBy=timers.target
